# systemd/iw4x-anticheat-updater.timer
[Unit]
Description=Check for stable IW4X anti-cheat updates every six hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
